# sub-hub/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── subconverter ──────────────────────────────────────────────────────────────
; Start first (priority=10) so that sub-server can reach it on boot.
[program:subconverter]
command=/usr/bin/env -u PORT /opt/subconverter/subconverter
directory=/opt/subconverter
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
; Redirect stdout/stderr to container's output
environment=HOME="/root"

; ── sub-server (Go) ───────────────────────────────────────────────────────────
; Start after subconverter (priority=20).
; sub-server has an internal wait loop that polls subconverter before serving.
[program:sub-server]
command=/usr/local/bin/sub-server
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
